# host/ansible/roles/coordinator/files/coordinator-throttle.service
[Unit]
# Record the Pi's undervoltage/throttle state to the journal for post-bounce
# forensics (#100). Fired by coordinator-throttle.timer.
Description=Record Pi power/throttle state to the journal

[Service]
Type=oneshot
ExecStart=/usr/local/bin/coord-throttle-log
